# supervisord.conf
[supervisord]
nodaemon=true

[program:magnolia]
command=/opt/magnolia/apache-tomcat/bin/magnolia_control.sh start
